# deployment_app/session_server/node/nodeapp.initrc
#!/sbin/openrc-run

name="nodeapp"
command="node /app/app.js -b"
pidfile="/run/${RC_SVCNAME}.pid"