# deployment_app/session_server/jamulus.initrc
#!/sbin/openrc-run

name="jamulus"
command="Jamulus -b"
pidfile="/run/${RC_SVCNAME}.pid"
command_args="-s -n"
output_log="/app/jamulus.log"